# archlinux/81-usb-pc-link.rules
ACTION=="add", ATTRS{idVendor}=="0471", ATTRS{idProduct}=="0111", RUN+="/usr/bin/usb_pc_link"
